# ugreen-truenas-zfs.timer
[Unit]
Description=Poll TrueNAS ZFS status for front-panel LEDs

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=ugreen-truenas-zfs.service

[Install]
WantedBy=timers.target
