# fan/ugreen-truenas-fan.service
[Unit]
Description=Drive UGREEN rear fan from TrueNAS disk temperatures
After=pve-guests.service
Conflicts=fancontrol.service

[Service]
Type=simple
EnvironmentFile=/etc/ugreen-truenas-fan.conf
ExecStart=/usr/local/bin/ugreen-truenas-fan.py --start
ExecStopPost=/usr/local/bin/ugreen-truenas-fan.py --stop
Restart=on-failure
RestartSec=10s
Nice=10

[Install]
WantedBy=multi-user.target
