# fan/ugreen-truenas-fan.timer
[Unit]
Description=Poll TrueNAS disk temperatures for UGREEN rear fan

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=ugreen-truenas-fan.service

[Install]
WantedBy=timers.target
